# 01-conf/runit/tomcat/finish
#!/bin/sh

/etc/tomcat-wrapper.sh stop >/dev/null 2>&1
