# 01-conf/runit/mongodb/run
#!/bin/sh

rm -f /var/lib/mongodb/mongod.lock
exec /usr/bin/mongod --config /etc/mongodb.conf --smallfiles --noprealloc
